# provision/app.service
[Unit]
Description=Gunicorn for Flask app
After=network.target

[Service]
User=www-data
Group=www-data
EnvironmentFile=/etc/myapp/app.env
WorkingDirectory=/opt/myapp
ExecStart=/usr/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app
Restart=always

[Install]
WantedBy=multi-user.target
